# rootfs/etc/systemd/system/sailtrack-lora2mqtt.service
[Unit]
Description=SailTrack-LoRa2MQTT

[Service]
ExecStart=/boot/sailtrack/sailtrack-lora2mqtt

[Install]
WantedBy=multi-user.target
