# rootfs/etc/systemd/system/sailtrack-kiosk.service
[Unit]
Description=SailTrack-Kiosk

[Service]
Restart=on-failure
ExecStart=/boot/sailtrack/sailtrack-kiosk

[Install]
WantedBy=multi-user.target
